# agents/runnerctl/contrib/systemd/runnerctl.service
[Unit]
Description=RunnerCTL API server (GitHub Actions self-hosted runner controller)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=runnerctl
Group=runnerctl
WorkingDirectory=/opt/runnerctl
Environment=PYTHONUNBUFFERED=1
# Bind API to loopback; expose via a tunnel if needed
Environment=RUNNERCTL_BIND=127.0.0.1:8080
# Load secrets and defaults
EnvironmentFile=-/etc/runnerctl/runnerctl.env
ExecStart=/usr/bin/python3 -m runnerctl.server
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target

